# docker/build/config/supervisor/queue.conf
[program:laravel-queue-default]
directory=/var/www
command=php artisan queue:work --tries=3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0